# deploy/systemd/awp-traffic-remote-backup.timer
[Unit]
Description=Daily AWP traffic monitor remote backup

[Timer]
OnCalendar=*-*-* 00:20:00
Persistent=true

[Install]
WantedBy=timers.target
